# target/linux/ramips/base-files/etc/hotplug.d/usb/30-usbledcheck
if [ "$ACTION" = "add" -o "$ACTION" = "remove" ];
then
	TRIGGER="$( cat /sys/class/leds/nw762b:green:usb/trigger | grep -o '\[.\+\]' )"
	[ $TRIGGER != "[usbport]" ] && echo "usbport" > /sys/class/leds/nw762b:green:usb/trigger
fi